# apps/api/signaldb-ai-dev-daemon.service
[Unit]
Description=SignalDB AI Dev Agent Daemon
Documentation=https://signaldb.live/docs
After=network.target postgresql.service
Wants=postgresql.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/signaldb/apps/api
EnvironmentFile=/etc/signaldb/api.env
ExecStart=/home/deploy/.bun/bin/bun run src/services/ai-dev-daemon.ts
Restart=on-failure
RestartSec=5

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signaldb-ai-dev-daemon

[Install]
WantedBy=multi-user.target
